# deploy/systemd/ollama-agent.service
[Unit]
Description=Ollama Agent Proxy Service
After=network.target
Documentation=https://github.com/yourusername/ollama-agent

[Service]
Type=simple
User=ollama
Group=ollama
ExecStart=/usr/local/bin/ollama-agent --local-addr 0.0.0.0:11434 --remote-url https://api.ollama.ai
Restart=on-failure
RestartSec=10
Environment=RUST_LOG=info
# Uncomment to specify an API key
#Environment=OLLAMA_API_KEY=your_api_key_here

# Security enhancements
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target